# scripts/systemd/c3po-gateway.service
[Unit]
Description=OpenClaw Gateway (C3PO family agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/ABSOLUTE/PATH/TO/THIS/REPO
EnvironmentFile=-%h/.config/c3po/.env
ExecStart=/usr/bin/env openclaw gateway --port ${GATEWAY_PORT:-18789} --auth-password ${GATEWAY_PASSWORD}
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=TZ=America/Sao_Paulo

[Install]
WantedBy=default.target
